# rebuild-minidlna.service
[Unit]
Description=Regenerate DLNA directories

[Service]
Type=simple
ExecStart=/opt/DLNA/dlna-einrichten.sh

[Install]
WantedBy=
